# 00-POC-24-Daemon/cris-daemon.service
[Unit]
Description=Cris Daemon

[Service]
Type=forking
ExecStart=/usr/bin/cris-daemon

[Install]
WantedBy=multi-user.target
